# deployment/ipl-fantasy.service
# Systemd service file for IPL Fantasy League
# Copy to /etc/systemd/system/ipl-fantasy.service

[Unit]
Description=IPL Fantasy Premier League 2026 Web Server
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/path/to/app
Environment="SECRET_KEY=your_secret_key_here"
Environment="RAPIDAPI_KEY=your_rapidapi_key_here"
Environment="ENVIRONMENT=production"
Environment="HOST=127.0.0.1"
Environment="PORT=8000"
Environment="RELOAD=false"
ExecStart=/path/to/app/venv/bin/python /path/to/app/app.py
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/path/to/app/data /path/to/app/Match data

[Install]
WantedBy=multi-user.target
